# deploy/aurhub-indexer.timer
[Unit]
Description=Hourly aurhub snapshot indexer

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
